# app/units/run-archive-reports.timer
[Unit]
Description=Run the proposed migration script (timer)

[Timer]
OnCalendar=*-*-* *:07,37:00

[Install]
WantedBy=timers.target
